# crontab
*/30 * * * * bundle exec rake /rails/lib/tasks/remind_task.rake remind_today_task
